# helpdesk/deploy/systemd/helpdesk-health.service
[Unit]
Description=Africa CDC Helpdesk — API health check
Documentation=file:///var/www/staff/helpdesk/documentation/SYSTEMD.md
PartOf=helpdesk.target
After=network-online.target
# After=apache2.service
# After=httpd.service

[Service]
Type=oneshot
User=@HELPDESK_USER@
Group=@HELPDESK_GROUP@
ExecStart=@HELPDESK_DEPLOY_BIN@/helpdesk-health.sh
