# apm/laravel12-queue-apm.service
[Unit]
Description=Laravel 12 Queue Worker for Africa CDC APM
After=network.target

[Service]
Type=simple
User=andrew
Group=www-data
WorkingDirectory=/var/lib/ACDC_SYSTEMS/staff/apm
ExecStart=/usr/bin/php artisan queue:work database --sleep=3 --tries=3 --timeout=90 --verbose
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=laravel-queue-apm

# Environment variables
Environment=APP_ENV=production
Environment=APP_DEBUG=false
Environment=QUEUE_CONNECTION=database

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
